# conf/supervisor.ini.example
[program:{{ project_name }}]
; todo write
